# supervisord-nginx.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/app/pb_logs/nginx.log
stderr_logfile=/app/pb_logs/nginx.log
priority=5

[program:pocketbase]
command=/app/pocketbase serve --http=0.0.0.0:8090 --dir=/app/pb_data --hooksDir=/app/pb_hooks --publicDir=/app/pb_public --migrationsDir=/app/pb_migrations
autostart=true
autorestart=true
stdout_logfile=/app/pb_logs/pocketbase.log
stderr_logfile=/app/pb_logs/pocketbase.log
priority=10

[program:nodejs]
command=node opencode-service.js
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/pb_logs/nodejs.log
stderr_logfile=/app/pb_logs/nodejs.log
priority=20

[program:healthcheck]
command=/app/healthcheck.sh --daemon
directory=/app
autostart=true
autorestart=true
stdout_logfile=/app/pb_logs/healthcheck.log
stderr_logfile=/app/pb_logs/healthcheck.log
priority=30